# deploy/snc-web.service
[Unit]
Description=S/NC Web SSR Server
After=snc-api.service

[Service]
Type=simple
User=snc
Group=snc
WorkingDirectory=/opt/snc/platform/apps/web
EnvironmentFile=/opt/snc/platform/.env
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=on-failure
RestartSec=5

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
