# deploy/snc-api.service
[Unit]
Description=S/NC API Server
After=network.target

[Service]
Type=simple
User=snc
Group=snc
WorkingDirectory=/opt/snc/platform/apps/api
EnvironmentFile=/opt/snc/platform/.env
ExecStart=/usr/bin/node --import tsx src/index.ts
Restart=on-failure
RestartSec=5


[Install]
WantedBy=multi-user.target
